# systemd/mytimer.service
[Unit]
Description=Run mytimer service

[Service]
ExecStart=/opt/mytimer.sh
